# Yuebon.WebApi/plytmsapi.ini
[program:plytmsapi]
command=/root/dotnet/dotnet Yuebon.WebApi.dll  --urls="http://*:5009" #要执行的命令
directory=/wwwroot/pyl.ts.yuebon.com #命令执行的目录
environment=ASPNETCORE__ENVIRONMENT=Production #环境变量
user=root  #进程执行的用户身份
stopsignal=INT
autostart=true #是否自动启动
autorestart=true #是否自动重启
startsecs=6 #自动重启间隔
stderr_logfile=/var/log/pyl.ts.yuebon.com.err.log #标准错误日志
stdout_logfile=/var/log/pyl.ts.yuebon.com.out.log #标准输出日志
